# deployment/systemd/communitas.service
[Unit]
Description=Communitas Node
After=network-online.target
Wants=network-online.target

[Service]
User=communitas
EnvironmentFile=-/etc/communitas-port.env
Environment=RUST_LOG=info,communitas=debug,saorsa_core=info
ExecStart=/opt/communitas/bin/communitas-headless --config /etc/communitas/config.toml --listen 0.0.0.0:${COMMUNITAS_QUIC_PORT} --metrics --metrics-addr 0.0.0.0:9600
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
