# deployment/systemd/communitas-updater.service
[Unit]
Description=Communitas Auto-Updater

[Service]
User=communitas
ExecStart=/opt/communitas/bin/communitas-autoupdater --config /etc/communitas/update.toml
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target